# systemd/clawbridge.service.tmpl
[Unit]
Description=ClawBridge Agent ({{SLUG}})
After=network.target

[Service]
Type=simple
ExecStart={{NODE_PATH}} {{PROJECT_ROOT}}/dist/index.js
WorkingDirectory={{HOME}}/.clawbridge
Restart=always
RestartSec=5
StandardOutput=append:{{HOME}}/.clawbridge/logs/out.log
StandardError=append:{{HOME}}/.clawbridge/logs/err.log
Environment="HOME={{HOME}}"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/snap/bin"
Environment="CONTAINER_IMAGE_BASE={{CONTAINER_IMAGE_BASE}}"
Environment="CONTAINER_IMAGE={{CONTAINER_IMAGE}}"
EnvironmentFile=-{{HOME}}/.clawbridge/.env

[Install]
WantedBy=default.target
